# example-images/atomic/cron.d/update-time
* * * * * root atomic run --spc larsks/atomic-persistent /atomic/bin/update-time > /dev/null 2>&1
